# workspace/crontab/laradock.platreform.example
* * * * * laradock /usr/bin/php /var/www/backend/artisan schedule:run >> /dev/null 2>&1
* * * * * laradock /usr/bin/php /var/www/landingpages/artisan schedule:run >> /dev/null 2>&1
