# workspace/crontab/laradock
* * * * * laradock /usr/bin/php /var/www/vipany-backend/artisan schedule:run >> /dev/null 2>&1
* * * * * laradock echo 'debug' > /var/www/cron.log >> /dev/null 2>&1